# scripts/com.sygen.host-update-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sygen.host-update-runner</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>__SCRIPT__</string>
        <string>--interval</string>
        <string>5</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>SYGEN_ROOT</key>
        <string>__SYGEN_ROOT__</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__SYGEN_ROOT__/logs/host-update-runner.log</string>
    <key>StandardErrorPath</key>
    <string>__SYGEN_ROOT__/logs/host-update-runner.log</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
