# scripts/com.sygen.keychain-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sygen.keychain-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__SCRIPT__</string>
        <string>--target</string>
        <string>__TARGET__</string>
        <string>--interval</string>
        <string>60</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__SYGEN_ROOT__/logs/keychain-sync.log</string>
    <key>StandardErrorPath</key>
    <string>__SYGEN_ROOT__/logs/keychain-sync.log</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
