# scripts/pro.sygen.updater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Sygen updater LaunchAgent (native macOS install). Polls GitHub Releases
  for new sygen / sygen-admin tags and exposes a 127.0.0.1:8082 HTTP
  endpoint that the core /apply flow calls to swap in a new venv/admin.

  Placeholders substituted by install.sh:
    __VENV_DIR__            absolute path to the Python venv
    __SYGEN_ROOT__          data/.env/logs root
    __SYGEN_UPDATER_PORT__  port LISTEN_PORT binds to (default 8082;
                            multi-instance auto-shifts when occupied)
    __HOME__                user's home dir
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>            <string>pro.sygen.updater</string>
    <key>ProgramArguments</key>
    <array>
        <string>__VENV_DIR__/bin/sygen-updater</string>
    </array>
    <key>WorkingDirectory</key> <string>__SYGEN_ROOT__</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SYGEN_HOME</key>          <string>__SYGEN_ROOT__/data</string>
        <key>SYGEN_UPDATER_PORT</key>  <string>__SYGEN_UPDATER_PORT__</string>
        <key>HOME</key>                <string>__HOME__</string>
        <key>PATH</key>                <string>__VENV_DIR__/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>        <true/>
    <key>KeepAlive</key>        <true/>
    <key>ThrottleInterval</key> <integer>10</integer>
    <key>StandardOutPath</key>  <string>__SYGEN_ROOT__/logs/updater.log</string>
    <key>StandardErrorPath</key><string>__SYGEN_ROOT__/logs/updater.log</string>
</dict>
</plist>
