# scripts/pro.sygen.admin.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Sygen admin LaunchAgent (native macOS install). Serves the Next.js
  standalone build at __ADMIN_DIR__/server.js on port __SYGEN_ADMIN_PORT__.

  Placeholders substituted by install.sh:
    __NODE_BIN__           absolute path to node (brew node@22)
    __ADMIN_DIR__          path to extracted admin tarball
    __SYGEN_ADMIN_PORT__   host port for admin (default 8080)
    __SYGEN_ROOT__         data/logs root
    __HOME__               user's home dir
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>            <string>pro.sygen.admin</string>
    <key>ProgramArguments</key>
    <array>
        <string>__NODE_BIN__</string>
        <string>__ADMIN_DIR__/server.js</string>
    </array>
    <key>WorkingDirectory</key> <string>__ADMIN_DIR__</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>             <string>production</string>
        <key>NEXT_TELEMETRY_DISABLED</key><string>1</string>
        <key>HOSTNAME</key>             <string>127.0.0.1</string>
        <key>PORT</key>                 <string>__SYGEN_ADMIN_PORT__</string>
        <key>HOME</key>                 <string>__HOME__</string>
    </dict>
    <key>RunAtLoad</key>        <true/>
    <key>KeepAlive</key>        <true/>
    <key>ThrottleInterval</key> <integer>10</integer>
    <key>StandardOutPath</key>  <string>__SYGEN_ROOT__/logs/admin.log</string>
    <key>StandardErrorPath</key><string>__SYGEN_ROOT__/logs/admin.log</string>
</dict>
</plist>
