# scripts/sygen-admin.service
[Unit]
Description=Sygen admin (Next.js standalone)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__ADMIN_DIR__
Environment=NODE_ENV=production
Environment=NEXT_TELEMETRY_DISABLED=1
Environment=HOSTNAME=127.0.0.1
Environment=PORT=__SYGEN_ADMIN_PORT__
ExecStart=__NODE_BIN__ __ADMIN_DIR__/server.js
Restart=always
RestartSec=5
StandardOutput=append:__SYGEN_ROOT__/logs/admin.log
StandardError=append:__SYGEN_ROOT__/logs/admin.log

[Install]
WantedBy=multi-user.target
