# scripts/sygen-compose.service
[Unit]
Description=Sygen Docker compose stack (auto-start on boot)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=__SYGEN_ROOT__
ExecStart=__DOCKER_BIN__ compose up -d
ExecStop=__DOCKER_BIN__ compose down
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
